# appdaemon/rootfs/etc/s6-overlay/s6-rc.d/appdaemon/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: AppDaemon
# Runs the AppDaemon
# ==============================================================================
declare log_level
declare debug_port
declare -a appdaemon_cmd
declare -a debug_cmd

bashio::log.info "Starting AppDaemon..."

# Find the matching Tor log level
log_level="INFO"
if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace|debug)
            log_level="DEBUG"
            ;;
        info|notice)
            log_level="INFO"
            ;;
        warning)
            log_level="WARNING"
            ;;
        error)
            log_level="ERROR"
            ;;
        fatal|off)
            log_level="FATAL"
            ;;
    esac
fi

# Run the AppDaemon
appdaemon_cmd=(appdaemon -c /config -D "${log_level}")

if bashio::config.true 'debug_enabled'; then
    debug_port="$(bashio::config 'debug_port')"
    debug_cmd=(python3 -m debugpy --listen "0.0.0.0:${debug_port}")

    if bashio::config.true 'debug_wait'; then
        debug_cmd+=(--wait-for-client)
    fi

    bashio::log.warning "Remote debugging enabled on port ${debug_port}"
    exec "${debug_cmd[@]}" "${appdaemon_cmd[@]}"
fi

exec "${appdaemon_cmd[@]}"
